# building_containers/building_app_container/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"

[program:node]
command=node app.js
